# systemd/.config/systemd/user/battery-notify.service
[Unit]
Description=Battery low notifications
After=graphical-session.target

[Service]
ExecStart=%h/.local/share/my-bins/battery-notify
Restart=always
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
